# .config/sys/etc/profile.d/30shinit.sh
if [ -r ~/.shinit ]; then
  . ~/.shinit
fi
