# .config/sys/etc/profile.d/30shinit.sh
if [ -r $HOME/.ashrc ]; then
  ENV=$HOME/.ashrc; export ENV
  . $ENV
fi
